# applications/luci-mcpd/root/uci-defaults/mcpd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mcpd[-1]
	add ucitrack mcpd
	set ucitrack.@mcpd[-1].init=mcpd
	commit ucitrack
EOF
rm -f /tmp/luci-indexcache
exit 0
