# themes/bootstrap/root/etc/uci-defaults/luci-theme-bootstrap
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Inteno=/luci-static/bootstrap
	commit luci
EOF
exit 0
