# applications/luci-dect/root/etc/uci-defaults/luci-dect
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dect[-1]
	add ucitrack dect
	set ucitrack.@dect[-1].init=dect
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
