# applications/luci-multiwan/root/etc/uci-defaults/luci-multiwan
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@multiwan[-1]
	add ucitrack multiwan
	set ucitrack.@multiwan[-1].init=multiwan
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
