# applications/luci-minidlna/root/etc/uci-defaults/luci-minidlna
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@minidlna[-1]
	add ucitrack minidlna
	set ucitrack.@minidlna[-1].init=minidlna
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
